# conf/etc/services.d/01-mattermost/run
#!/bin/sh
NAME="MATTERMOST"
COMMAND="source /tmp/mattermost_env; /mattermost/bin/mattermost server"
WAIT="/tmp/mattermost"

exec service-handler "$NAME" "$COMMAND" "$WAIT"
